# services/trader-daemon.service
[Unit]
Description=Autonomous Oil Trader — Order/State Daemon
After=trader-stream.service network-online.target

[Service]
Type=simple
WorkingDirectory=%h/autonomous-oiltrader
EnvironmentFile=%h/autonomous-oiltrader/.env
ExecStart=/usr/bin/env python -m daemons.trader_daemon
Restart=on-failure
RestartSec=30
StandardOutput=append:%h/autonomous-oiltrader/logs/trader.log
StandardError=append:%h/autonomous-oiltrader/logs/trader.err

[Install]
WantedBy=default.target
