# services/trader-api.service
[Unit]
Description=Autonomous Oil Trader — FastAPI Server
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/autonomous-oiltrader
EnvironmentFile=%h/autonomous-oiltrader/.env
ExecStart=/usr/bin/env uvicorn interfaces.api:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=10
StandardOutput=append:%h/autonomous-oiltrader/logs/api.log
StandardError=append:%h/autonomous-oiltrader/logs/api.err

[Install]
WantedBy=default.target
